# pacicron.service
[Unit]
Description=pacicron - An pacman update notifier

[Service]
Type=oneshot
ExecStart=/usr/bin/pacicron.sh
